# recipes-netkit/netkit-telnet/netkit-telnet/telnet.socket
[Unit]
Description=Telnet Server Activation Socket
Documentation=man:telnetd(8)

[Socket]
ListenStream=23
Accept=yes

[Install]
WantedBy=sockets.target
